# 7.1/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=true

[program:php-fpm]
command=php-fpm -F
user=root
autostart=true

[supervisorctl]